# docker-files/dev/config/etc/services.d/laraswoole/run
#!/bin/sh

cd /var/www/html

php artisan swoole:http start